# buildroot/overlay/lib/udev/rules.d/61-usbdevice.rules
SUBSYSTEM=="usb", ACTION=="add", RUN+="/usr/bin/usbdevice"
SUBSYSTEM=="block", ACTION=="add", RUN+="/usr/bin/sd-cards-auto-mount"
